# src/myautoconnectscript.conf
start on runlevel [2345]
stop on runlevel [!2345]
exec /home/pi/dev/rounter-RUC/src/auto-connect.py
